# Livebox 3 Pro/firmware/etc/init.d/ppp
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 Sagemcom
START=550

start() {
	mkdir -p /tmp/etc/ppp
}

