# Livebox 4/etc/profile.d/tracker.sh
export TRACKER_USE_CONFIG_FILES=1
export TRACKER_BUS_TYPE=SYSTEM

