# Livebox 4/lib/udev/rules.d/75-storage.rules
SUBSYSTEM=="block", KERNEL=="sd*", RUN+="/usr/lib/storage/storage-udev.sh %k"
SUBSYSTEM=="usb", RUN+="/usr/lib/storage/storage-udev.sh %k"
