# Livebox 4/etc/inittab
# disaster recovery if necessary
::wait:/etc/init.d/sysinit

# shell on serial port
::askfirst:-/bin/login

# shutdown scripts to run
::shutdown:/etc/init.d/rc 6 K stop

