# Livebox 3 Pro/firmware/etc/rc.d/S22orange-passwd
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 Sagemcom
START=22

start() {
    echo -n "Init LiveBox password: "
	orange-passwd
    echo "OK"
}

