# Livebox 2/sagem/rescue/fs/rw/etc/inittab
# disaster recovery if necessary
::wait:/etc_ro_fs/init.d/sysinit

# system initialisation
::wait:/etc_ro_fs/init.d/rc.sysinit

# bash
::askfirst:-/bin/login
#::askfirst:-/bin/sh --login
